# crontab
* * * * * cd /app && python3 time_to_party.pyw
#0 13 * * * cd /app && python3 time_to_party.pyw
